# supervisor/roles/files/supervisord_tomcat-8001-boss.conf
[program:tomcat-8001-boss]
numprocs = 1
process_name = tomcat-8001-boss
directory = /services/module/tomcat-8001-boss
environment=JAVA_HOME="/usr/local/jdk1.8.0_45/",JAVA_BIN="/usr/local/jdk1.8.0_45/bin/"
command = /services/module/tomcat-8001-boss/bin/catalina.sh run
autostart = true 
autorestart = true
stdout_logfile = /var/log/supervisord_tomcat-8001-boss_acess.log
stderr_logfile = /var/log/supervisord_tomcat-8001-boss_error.log
stdout_logfile_maxbytes = 256MB
stdout_logfile_backups = 31
